# files/v2raya-policy-boot
#!/bin/sh /etc/rc.common
START=99
STOP=10

start() {
  (
    sleep 12
    /etc/init.d/v2raya enable >/dev/null 2>&1 || true
    /etc/init.d/v2raya restart >/dev/null 2>&1 || /etc/init.d/v2raya start >/dev/null 2>&1 || true
    /etc/init.d/uhttpd enable >/dev/null 2>&1 || true
    /etc/init.d/uhttpd restart >/dev/null 2>&1 || /etc/init.d/uhttpd start >/dev/null 2>&1 || true
    /usr/bin/v2raya-policy-apply >/tmp/v2raya-policy-boot.log 2>&1 || true
    /usr/bin/v2raya-device-policy >>/tmp/v2raya-policy-boot.log 2>&1 || true
    /usr/bin/v2raya-dns-policy >>/tmp/v2raya-policy-boot.log 2>&1 || true
  ) &
}
